# affinity.service
[Unit]
Description=Apply minimal display priority tuning

# Run after display is up so the script can tune display.service safely.
After=display.service

# Re-run when display.service is restarted/stopped.
PartOf=display.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/mks/display_connector/affinity-setup.sh

[Install]
WantedBy=display.service
